# tests/docker_postgresql_bucardo_alpine/etc/supervisord.conf
[supervisord]
logfile = /data/logs/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = /data/pids/supervisord.pid
nodaemon = true
minfds = 1024
minprocs = 200
umask = 022
identifier = supervisor
directory = /data
nocleanup = true
childlogdir = /data/logs
strip_ansi = false

[unix_http_server]
file = /data/pids/supervisor.sock

[supervisorctl]
serverurl = unix:///data/pids/supervisor.sock     

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:precheck.sh]
command=/usr/local/sbin/precheck.sh
user=root
priority=100
autorestart=false
startretries=0
startsec=60

[program:postgresql]
command=/usr/bin/postgres -D /data/postgresql -c config_file=/data/postgresql/postgresql.conf
user=postgres
autostart=false
autorestart=true
priority=200

[program:postcheck.py]
command=/usr/local/sbin/postcheck.py
user=root
priority=999
autostart=false
autorestart=true
startretries=3

[program:zabbix_agentd]
command = /usr/sbin/zabbix_agentd -c /etc/zabbix/zabbix_agentd.conf -f
user=zabbix
startsecs = 30
autorestart=true
priority=300
